# ansible/physical/roles/btrfs-balance/templates/btrfs-balance.timer
[Unit]
Description=Periodically start a btrfs balance

[Timer]
OnCalendar=Monthly
Persistent=True

[Install]
WantedBy=timers.target
