# ansible/roles/reflector/files/etc_systemd_system/reflector-update.timer
[Unit]
Description=Periodically update Pacman's mirrorlist

[Timer]
OnCalendar=Mon,Wed,Fri
Persistent=true

[Install]
WantedBy=timers.target
