# ansible/physical/roles/borg-backup-users/templates/borg-backup.service
[Unit]
Description=Back up files with Borg

[Service]
Type=oneshot
ExecStart={{ borg_backup_user_home }}/.local/bin/borg-backup.sh
