# ansible/physical/roles/aur/files/config_systemd_user/update-packages.timer
[Unit]
Description=Periodically update AUR packages

[Timer]
OnCalendar=daily
Persistent=True

[Install]
WantedBy=timers.target
