# ansible/roles/subsonic/templates/subsonic-backup.timer.j2
[Unit]
Description=Periodically back up subsonic's database rebuild script

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
