# ansible/roles/btrfs-scrub/templates/btrfs-scrub.timer.j2
[Unit]
Description=Periodically start a btrfs scrub

[Timer]
OnCalendar=Monthly
Persistent=True

[Install]
WantedBy=timers.target
