# ansible/roles/btrfs-scrub/templates/btrfs-scrub.service.j2
[Unit]
Description=Start a btrfs scrub

[Service]
ExecStart=/usr/bin/btrfs scrub start -Bd '{{ btrfs_path }}'
ExecStop=/usr/bin/btrfs scrub cancel '{{ btrfs_path }}'
