# ansible/physical/roles/smartd/files/etc_systemd_system/smart-test-short.timer
[Unit]
Description=Start a short SMART scan every Wed and Fri

[Timer]
OnCalendar=Wed,Fri
Persistent=true

[Install]
WantedBy=timers.target
