# ansible/physical/roles/smartd/files/etc_systemd_system/smart-test-long.timer
[Unit]
Description=Start a long SMART scan every Mon

[Timer]
OnCalendar=Mon
Persistent=true

[Install]
WantedBy=timers.target
