# ansible/physical/roles/smartd/files/etc_systemd_system/smart-test-long.service
[Unit]
Description=Start a long SMART scan

[Service]
Type=oneshot
ExecStart=/usr/local/bin/smart-test-long
