# ansible/physical/roles/smartd/files/etc_systemd_system/smart-test-short.service
[Unit]
Description=Start a short SMART scan

[Service]
Type=oneshot
ExecStart=/usr/local/bin/smart-test-short
