# ansible/physical/roles/dynamic-dns/files/etc_systemd_system/dynamic-dns.service
[Unit]
Description=Update dynamic DNS

[Service]
Type=oneshot
ExecStart=/usr/local/bin/dynamic-dns
