# ansible/roles/borg-backup-ichimonji10/files/config_systemd_user/borg-backup.service
[Unit]
Description=Back up files with Borg

[Service]
Type=oneshot
ExecStart=/home/ichimonji10/.local/bin/borg-backup
