# ansible/roles/hostapd/files/log-entropy.timer
[Unit]
Description=Periodically log available entropy

[Timer]
OnBootSec=120
OnUnitActiveSec=300

[Install]
WantedBy=timers.target
