# ansible/roles/btrfs-scrub/files/etc_systemd_system/btrfs-scrub.service
[Unit]
Description=Start a btrfs scrub

[Service]
Type=oneshot
ExecStart=/usr/bin/btrfs scrub start /mnt/btrfs
