# ansible/roles/btrfs-scrub/files/etc_systemd_system/btrfs-scrub.service
[Unit]
Description=Start a btrfs scrub

[Service]
ExecStart=/usr/bin/btrfs scrub start -Bd /mnt/btrfs
