# ansible/physical/roles/mnemosyne/templates/mnemosyne-sync-server.service
[Unit]
Description=Mnemosyne sync server

[Service]
ExecStart=/usr/bin/mnemosyne --sync-server
User={{ mnemosyne_user }}
Group={{ mnemosyne_user }}
