# ansible/physical/roles/libvirtd/files/etc_systemd_system/virtlogd.service.d/override.conf
# Execute `systemctl cat virtlogd` to see the complete unit.
[Unit]
StopWhenUnneeded=true
