# ansible/physical/roles/hostapd/files/hostapd.service.d/override.conf
[Unit]
Wants=wlan-secure.service
Before=wlan-secure.service
