# ansible/physical/roles/mnemosyne/files/sync-server.service
[Unit]
Description=Mnemosyne sync server

[Service]
ExecStart=/usr/bin/mnemosyne --sync-server
