# ansible/physical/roles/libvirtd/files/etc_systemd_system/dnsmasq.service.d/override.conf
# Execute `systemctl cat dnsmasq` to see the complete unit.
[Unit]
StopWhenUnneeded=true
